# conf/refresh-site.cron
*/10 * * * * @USER@ /srv/www.rabbitmq.com/code/refresh-site
