# conf/sync-nightly.cron
0 * * * * @USER@ /srv/www.rabbitmq.com/code/sync-nightly
